# memory_usage.service
[Unit]
Description=Memory usage service

[Service]
ExecStart=/bin/bash /usr/sbin/memory_usage.py

[Install]
WantedBy=multi-user.target

